# roles/extractor/scripts/systemd/iwe-extractor-inbox-check.timer
[Unit]
Description=IWE Extractor — inbox-check 3 раза в день (09:00, 15:00, 21:00 UTC)

[Timer]
OnCalendar=*-*-* 09:00:00
OnCalendar=*-*-* 15:00:00
OnCalendar=*-*-* 21:00:00
Persistent=true
Unit=iwe-extractor-inbox-check.service

[Install]
WantedBy=timers.target
